# daemon/openrc/nyansd
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/usr/local/bin/nyansd"
command_args=""
